# files/etc/rc.local
uci set network.wwan=interface
uci set network.wwan.proto=static
uci set network.wwan.ipaddr=192.168.100.100
uci set network.wwan.gateway=192.168.1.1
uci set network.wwan.netmask=255.255.255.0
uci set network.wwan.dns=8.8.8.8
uci set network.lan.ipaddr=192.168.1.100
uci set network.lan.gateway=192.168.1.1
uci set network.lan.dns=8.8.8.8
uci commit network
/etc/init.d/uhttpd enable
/etc/init.d/uhttpd start
/etc/init.d/network restart
/etc/enablestart
exit 0